# backend/crontab
# Scheduled data updates for Tapas FPL App
#
# Run daily at 08:00 UTC (08:00 GMT / 09:00 BST)
# After most GWs finalize (~22:00 UK time) + buffer for bonus points
# Changed from 06:00 to 08:00 to ensure FPL has set data_checked flag
#
# Updates:
# - Points Against: FPL points conceded by each team (~5-10 min incremental)
# - Manager Snapshots: GW history and picks for H2H comparison (~1-2 min)
# - Chips Usage: Manager chip activations for tracked league (~30 sec)
#
# Idempotent - safe to run multiple times (skips if no new finalized GW)

0 8 * * * python -m scripts.scheduled_update
