# backend/crontab
# Scheduled data updates for Tapas FPL App
#
# Run daily at 06:00 UTC (06:00 GMT / 07:00 BST)
# After most GWs finalize (~22:00 UK time) + buffer for bonus points
#
# Updates:
# - Points Against: FPL points conceded by each team (~5-10 min incremental)
# - Manager Snapshots: GW history and picks for H2H comparison (~1-2 min)
# - Chips Usage: Manager chip activations for tracked league (~30 sec)
#
# Idempotent - safe to run multiple times (skips if no new finalized GW)

0 6 * * * python -m scripts.scheduled_update
